# scripts/systemd/goio-metrics.timer
[Unit]
Description=GOIO metrics timer (hourly)

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
