# palacio-central/scripts/systemd/goio-metrics.service
[Unit]
Description=GOIO host+docker metrics snapshot
Documentation=https://github.com/<org>/<repo>
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc '/usr/local/bin/collect_metrics.sh'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUSchedulingPolicy=idle
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/log/goio-metrics
